# ops/systemd/asdev-agent-loop.timer
[Unit]
Description=ASDEV Autonomous Agent Loop Timer
After=default.target network-online.target

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
RandomizedDelaySec=2min
Unit=asdev-agent-loop.service

[Install]
WantedBy=timers.target
