# ops/systemd/asdev-agent-loop.service
[Unit]
Description=ASDEV Autonomous Agent Loop
After=default.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/dev13/my-project/sites/live/alirezasafaeisystems
Environment=HOME=/home/dev13
Environment=PATH=/home/dev13/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=ASDEV_COMMAND_ISSUE=45
Environment=ASDEV_MAX_JOBS=2
Environment=ASDEV_AGENT_LOG_DIR=/tmp/asdev-agent-loop/logs
ExecStart=/usr/bin/env bash -lc '/home/dev13/my-project/scripts/agent-command-center/run-autonomous-loop.sh --issue 45 --max-jobs 2'
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asdev-agent-loop
