# ops/systemd/asdev-agent-loop.timer
[Unit]
Description=ASDEV Autonomous Agent Loop Timer (every 30m)
After=default.target network-online.target

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Persistent=true
RandomizedDelaySec=90
Unit=asdev-agent-loop.service

[Install]
WantedBy=timers.target
