# ops/systemd/monitoring/asdev-monitor-public-http.timer
[Unit]
Description=ASDEV public HTTP probe every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
Unit=asdev-monitor-public-http.service

[Install]
WantedBy=timers.target
